# init/gtc-configurator.conf
# gtc-configurator
#
# gtc-configurator configure network and vmware

description     "gtc-configurator"

start on runlevel [35]

exec gtc-configurator

