# init/gtc-update.conf
# gtc-configurator
#
# gtc-configurator configure network and vmware

description     "gtc-configurator"

start on runlevel [2345]
stop on runlevel [!2345]

expect fork
respawn

exec gtc-configurator

