# init/gtc-update.conf
# gtc-configurator
#
# gtc-configurator configure network and vmware

description     "gtc-update"

start on networking


script
gtc-update -c
end script
